# deploy/vermilion-bird.service
[Unit]
Description=Vermilion Bird - LLM Chat Service
Documentation=https://github.com/your-org/vermilion-bird
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vermilion-bird
Group=vermilion-bird

# 工作目录
WorkingDirectory=/opt/vermilion-bird

# 环境变量
Environment="PATH=/opt/vermilion-bird/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="LANG=en_US.UTF-8"
Environment="LC_ALL=en_US.UTF-8"

# 可选：通过环境文件加载敏感配置
# EnvironmentFile=/etc/vermilion-bird/config.env

# 启动命令 - 飞书服务模式
# 如需其他模式，修改 ExecStart 即可
ExecStart=/opt/vermilion-bird/.venv/bin/vermilion-bird feishu

# 重启策略
Restart=always
RestartSec=10

# 资源限制
LimitNOFILE=65535
# MemoryMax=1G
# CPUQuota=50%

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vermilion-bird

# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/vermilion-bird /var/lib/vermilion-bird /var/log/vermilion-bird
ReadOnlyPaths=/etc

# 优雅关闭
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
